# deploy/com.berat.btcbot.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.berat.btcbot.dashboard</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/beratdal/btcbot_runtime/deploy/run_dashboard.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/beratdal/btcbot_runtime</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/beratdal/btcbot_runtime/logs/dashboard.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/beratdal/btcbot_runtime/logs/dashboard.err</string>
  </dict>
</plist>
